# rpi/physicalidd.service
[Unit]
Description=Physical ID Daemon
After=network.target

[Service]
ExecStart=/usr/bin/physicalidd
Restart=always
User=root


[Install]
WantedBy=multi-user.target
